# system/systemd/user/moltbook-relay.service
[Unit]
Description=Relay moltbook staging files to sanitizer input

[Service]
Type=oneshot
ExecStart={{ HOME }}/.openclaw/workspace/scripts/moltbook-relay.sh
