# system/systemd/user/openclaw-gateway.service.d/10-moltbook-env.conf
[Service]
EnvironmentFile={{ HOME }}/.openclaw/moltbook.env
